# app/deploy/systemd/apetitas-meta-nightly.service
[Unit]
Description=Apetitas - meta nightly (fill missing meta_title/meta_description)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=deploy
WorkingDirectory=/home/deploy/backend/app
ExecStart=/home/deploy/backend/app/.venv/bin/python /home/deploy/backend/app/manage.py run_recipe_meta_nightly --limit=500
# Log to journald (view via: journalctl -u apetitas-meta-nightly.service)

[Install]
WantedBy=multi-user.target
