# app/deploy/systemd/apetitas-meta-nightly.timer
[Unit]
Description=Apetitas - meta nightly timer

[Timer]
OnCalendar=*-*-* 02:20:00
Persistent=true

[Install]
WantedBy=timers.target
