# app/deploy/systemd/apetitas-nutrition-poll.timer
[Unit]
Description=Apetitas - nutrition batch poll timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
