# app/deploy/systemd/apetitas-nutrition-nightly.timer
[Unit]
Description=Apetitas - nutrition nightly timer

[Timer]
OnCalendar=*-*-* 02:10:00
Persistent=true

[Install]
WantedBy=timers.target
